# deploy/habbig-world.service
[Unit]
Description=Habbig World State Dashboard (port 7050)
After=network.target

[Service]
# Required env vars: (none — world-state dashboard reads no env vars from .env)
# Note: This service loads the shared gateway .env but only needs the above vars
Type=simple
User=julianhabbig
WorkingDirectory=/home/julianhabbig/Polymarket/world-state-dashboard
ExecStart=/home/julianhabbig/Polymarket/venv/bin/python -m uvicorn server:app --host 127.0.0.1 --port 7050
Restart=always
RestartSec=5
EnvironmentFile=-/home/julianhabbig/Polymarket/gateway/.env.production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/julianhabbig/Polymarket/world-state-dashboard

[Install]
WantedBy=multi-user.target
