# deploy/habbig-stock.service
[Unit]
Description=Habbig Stock Dashboard (port 8050)
After=network.target

[Service]
Type=simple
User=julianhabbig
WorkingDirectory=/home/julianhabbig/Polymarket/stock-dashboard
ExecStart=/home/julianhabbig/Polymarket/venv/bin/python stock_dashboard.py --port 8050
Restart=always
RestartSec=5
EnvironmentFile=-/home/julianhabbig/Polymarket/gateway/.env.production

[Install]
WantedBy=multi-user.target
