# deploy/narve-stock.service
[Unit]
Description=Narve Stock Dashboard (port 8050)
After=network.target

[Service]
# Required env vars: (none — stock dashboard reads no env vars from .env)
# Note: This service loads the shared gateway .env but only needs the above vars
Type=simple
User=julianhabbig
WorkingDirectory=/home/julianhabbig/Polymarket/stock-dashboard
ExecStart=/home/julianhabbig/Polymarket/venv/bin/python stock_dashboard.py --port 8050
Restart=always
RestartSec=5
EnvironmentFile=-/home/julianhabbig/Polymarket/gateway/.env.production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/julianhabbig/Polymarket/stock-dashboard
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
