# deploy/narve-gateway.service
[Unit]
Description=Narve Gateway (port 7000)
After=network.target
Wants=narve-crypto.service narve-weather.service narve-sports.service narve-world.service narve-midterm.service narve-traders.service narve-stock.service narve-whale.service

[Service]
# Required env vars: GATEWAY_SSO_SECRET, GATEWAY_COOKIE_SECRET, STRIPE_SECRET_KEY,
#   STRIPE_PUBLISHABLE_KEY, STRIPE_WEBHOOK_SECRET, SMTP_USER, SMTP_PASS,
#   SMTP_HOST, SMTP_PORT, ENQUIRY_EMAIL, PRODUCTION
# Note: This service loads the shared gateway .env and uses most vars
Type=simple
User=julianhabbig
WorkingDirectory=/home/julianhabbig/Polymarket/gateway
ExecStart=/home/julianhabbig/Polymarket/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PRODUCTION=1
EnvironmentFile=-/home/julianhabbig/Polymarket/gateway/.env.production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/julianhabbig/Polymarket/gateway
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
