# deploy/narve-sports.service
[Unit]
Description=Narve Sports Dashboard (port 8888)
After=network.target

[Service]
# Required env vars: GATEWAY_SSO_SECRET, ODDS_API_KEY
# Note: This service loads the shared gateway .env but only needs the above vars
Type=simple
User=julianhabbig
WorkingDirectory=/home/julianhabbig/Polymarket/sports-dashboard
ExecStart=/home/julianhabbig/Polymarket/venv/bin/python sports_dashboard.py
Restart=always
RestartSec=5
EnvironmentFile=-/home/julianhabbig/Polymarket/gateway/.env.production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/julianhabbig/Polymarket/sports-dashboard
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
