# deploy/narve-traders.service
[Unit]
Description=Narve Top Traders Dashboard (port 8052)
After=network.target

[Service]
# Required env vars: (none — top-traders dashboard reads no env vars from .env)
# Note: This service loads the shared gateway .env but only needs the above vars
Type=simple
User=julianhabbig
WorkingDirectory=/home/julianhabbig/Polymarket/top-traders-dashboard
ExecStart=/home/julianhabbig/Polymarket/venv/bin/python server.py
Restart=always
RestartSec=5
EnvironmentFile=-/home/julianhabbig/Polymarket/gateway/.env.production

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/julianhabbig/Polymarket/top-traders-dashboard
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
